# deploy/supervisor/service.conf
[program:[service name]]

command=/path/to/virtualenv/bin/python /path/to/app/service/web.py
OR
command=gunicorn -c /path/to/gunicorn_config.py service.web:app

user=[user to run as]
directory=/path/to/app/
stdout_logfile=/var/log/supervisor/%(program_name)s-access.log
stderr_logfile=/var/log/supervisor/%(program_name)s-error.log
autostart=true
autorestart=true
stopasgroup=true